# .docker/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command=/usr/local/sbin/php-fpm -F
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/php-fpm.log
